# deploy/amber.service
# systemd unit for Amber — install on the OVH VPS at /etc/systemd/system/amber.service
#
# Assumes:
#   - app deployed to /opt/amber  (clone or rsync this repo there)
#   - a virtualenv at /opt/amber/.venv with deps installed
#   - /opt/amber/.env present (copied from .env.example, secrets filled in)
#   - a dedicated unprivileged user `amber`
#
# See deploy/README.md for the full first-time setup.

[Unit]
Description=Amber voice backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=amber
Group=amber
WorkingDirectory=/opt/amber
EnvironmentFile=/opt/amber/.env
ExecStart=/opt/amber/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3

# Hardening.
# NOTE: NoNewPrivileges=true blocks the voice self-update (update_server), because
# the update command has to `sudo` to restart the service. If you want hands-free
# backend updates, set this to `no` (update.sh's actual work runs in a separate,
# detached systemd-run unit, so the running service stays otherwise sandboxed).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# ProtectHome=true hides /home — keep the app OUT of /home (e.g. /opt/amber) or
# the service can't read its own code. Set to `no` if the checkout lives in a home dir.
ProtectHome=true

[Install]
WantedBy=multi-user.target
